# systemd/obsidian-mcp.service.template
[Unit]
Description=Obsidian MCP Server (Docker Compose)
After=network.target obsidian-sync.service
Requires=obsidian-sync.service
PartOf=obsidian-mcp.target

[Service]
Type=simple
WorkingDirectory=__PROJECT_DIR__
ExecStart=__DOCKER_BIN__ compose up --build
ExecStop=__DOCKER_BIN__ compose down
Restart=on-failure
RestartSec=5

[Install]
WantedBy=obsidian-mcp.target
