# systemd/obsidian-sync.service
[Unit]
Description=Obsidian Headless Sync
After=network.target
PartOf=obsidian-mcp.target

[Service]
Type=simple
EnvironmentFile=/root/obsidian-mcp/.env
ExecStart=/root/.local/share/pnpm/ob sync --continuous --path ${VAULT_PATH}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=obsidian-mcp.target
