# systemd/obsidian-mcp.target.template
[Unit]
Description=Obsidian MCP Stack
Wants=__SYNC_SERVICES__ obsidian-mcp.service caddy.service

[Install]
WantedBy=multi-user.target
