# systemd/obsidian-mcp.service
[Unit]
Description=Obsidian MCP Server (Docker Compose)
After=network.target obsidian-sync.service
Requires=obsidian-sync.service
PartOf=obsidian-mcp.target

[Service]
Type=simple
WorkingDirectory=/root/obsidian-mcp
ExecStart=/usr/bin/docker compose up --build
ExecStop=/usr/bin/docker compose down
Restart=on-failure
RestartSec=5

[Install]
WantedBy=obsidian-mcp.target
