# systemd/obsidian-mcp.target.template
[Unit]
Description=Obsidian MCP Stack
Wants=obsidian-sync.service obsidian-mcp.service

[Install]
WantedBy=multi-user.target
