# systemd/obsidian-mcp.target
[Unit]
Description=Obsidian MCP Stack (sync + server + caddy)

[Install]
WantedBy=multi-user.target
