# systemd/obsidian-sync.service
[Unit]
Description=Obsidian Headless Sync
After=network.target
PartOf=obsidian-mcp.target

[Service]
Type=simple
WorkingDirectory=/root/obsidian-mcp/vault
ExecStart=/root/.local/share/pnpm/ob sync --continuous
Restart=on-failure
RestartSec=5

[Install]
WantedBy=obsidian-mcp.target
