# systemd/obsidian-sync.service
[Unit]
Description=Obsidian Headless Sync
After=network.target
PartOf=obsidian-mcp.target

[Service]
Type=simple
ExecStart=/bin/bash -c 'cd "$VAULT_PATH" && exec /root/.local/share/pnpm/ob sync --continuous'
EnvironmentFile=/root/obsidian-mcp/.env
Restart=on-failure
RestartSec=5

[Install]
WantedBy=obsidian-mcp.target
